# deploy/sahayak.service
[Unit]

Description=Sahayak Robot Server

After=network-online.target

Wants=network-online.target



[Service]
# Adjust paths above/below to wherever you cloned the repository
Environment=PYTHONPATH=/app/pydev_demo/02_detection_sample/03_ultralytics_yolov8

Type=simple

User=sunrise

WorkingDirectory=/home/sunrise/sahayak-guardian-ai

ExecStart=/usr/bin/python3 /home/sunrise/sahayak-guardian-ai/src/server.py

Restart=on-failure

RestartSec=5

# Give the network/hotspot a moment to come up before starting

ExecStartPre=/bin/sleep 15

# Online mode (optional): uncomment and add your key later




[Install]

WantedBy=multi-user.target